# bedrock-server-backup.timer
[Unit]
Description=Bedrock Server Backup Service Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=15m

[Install]
WantedBy=timers.target
